# roles/nginx_letsencrypt/templates/letsencrypt-alert@.service.j2
[Unit]
Description=Service for alerting when timer unit for letsencrypt renewal fails
After=network.target

[Service]
Type=simple
ExecStart={{ letsencrypt_alert }}
