# roles/nginx/files/nginx.service.d/extend_file_limit.conf
# extend the limit of open file descriptors by nginx up from the default of 1024 on Ubuntu
#
# installed via ansible: do not modify here

[Service]
LimitNOFILE=15000
